# ops/slowlink-watchdog.service
[Unit]
Description=SlowLink App CPU Watchdog
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
Environment=APP_CONTAINER=slowlink_app
Environment=CHECK_INTERVAL=20
Environment=CPU_THRESHOLD=90
Environment=HIGH_COUNT_LIMIT=4
Environment=COOLDOWN_SECONDS=600
Environment=LOG_FILE=/opt/slowlink/watchdog.log
ExecStart=/bin/sh /opt/slowlink/ops/slowlink_watchdog.sh
Restart=always
RestartSec=5
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
